# src/package/luci-app-usb_otg/root/etc/uci-defaults/luci-app-usb_otg
#!/bin/sh
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@usb_otg[-1]
	add ucitrack usb_otg
	set ucitrack.@usb_otg[-1].init=usb_otg
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0

